# docker/supervisord.conf
; Marveen process supervisor — replaces systemd/launchd inside the container.
; Runs as the (non-root) marveen user. tini is PID 1 and reaps the orphan
; bun pollers that outlive `tmux kill-session`.

[supervisord]
nodaemon=true
loglevel=info
logfile=/app/store/supervisord.log
pidfile=/app/store/supervisord.pid
childlogdir=/app/store

[unix_http_server]
file=/app/store/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///app/store/supervisor.sock

; ── Dashboard: HTTP server + in-process pollers/watchers (the orchestrator) ──
; The real entrypoint is dist/index.js (package.json `start`). NOTE: start.sh's
; container fallback launches `bun src/web/serve.ts`, which DOES NOT EXIST — do
; not use it.
[program:dashboard]
command=node /app/dist/index.js
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=1000
stopsignal=TERM
stopwaitsecs=10
priority=10
stdout_logfile=/app/store/dashboard.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/app/store/dashboard.error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

; ── Channels: the main-agent Claude Code bridge in a tmux session ──
; channels.sh exits when its tmux session dies; autorestart respawns it
; (equivalent to systemd Restart=always). It has its own rapid-failure backoff.
[program:channels]
command=bash /app/scripts/channels.sh
directory=/app
autostart=true
autorestart=true
startsecs=20
startretries=1000
stopsignal=TERM
stopwaitsecs=15
priority=20
stdout_logfile=/app/store/channels.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/app/store/channels.error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
